# deployment/supervisor.ini
[program:karyaradhane]
command=/opt/karyaradhane/backend/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
directory=/opt/karyaradhane/backend
user=www
autostart=true
autorestart=true
stderr_logfile=/var/log/karyaradhane.err.log
stdout_logfile=/var/log/karyaradhane.out.log
environment=PYTHONPATH="/opt/karyaradhane/backend",PATH="/opt/karyaradhane/backend/venv/bin:/usr/local/bin:/usr/bin:/bin"
startsecs=10
stopwaitsecs=600
killasgroup=true
stopasgroup=true

[program:karyaradhane-worker]
; Add background worker if needed for async tasks
; command=/opt/karyaradhane/backend/venv/bin/python worker.py
; directory=/opt/karyaradhane/backend
; user=www
; autostart=false
; autorestart=true
; stderr_logfile=/var/log/karyaradhane-worker.err.log
; stdout_logfile=/var/log/karyaradhane-worker.out.log
